# udev/rules.d/70-snap.vlc.rules
# This file is automatically generated.
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_vlc_vlc"
# opengl
KERNEL=="vchiq", TAG+="snap_vlc_vlc"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_vlc_vlc"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_vlc_vlc"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_vlc_vlc"
TAG=="snap_vlc_vlc", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_vlc_vlc $devpath $major:$minor"
