# udev/rules.d/70-snap.spotify.rules
# This file is automatically generated.
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_spotify_spotify"
# opengl
KERNEL=="vchiq", TAG+="snap_spotify_spotify"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_spotify_spotify"
TAG=="snap_spotify_spotify", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_spotify_spotify $devpath $major:$minor"
